# deploy/anchord.service
[Unit]
Description=dreamtree anchor seam (HTTP -> MsgCommitSeed)
After=network-online.target dreamtreed.service
Wants=dreamtreed.service
# The account/mint key (keyring-backend test) lives ONLY on a TPM-unlocked
# LUKS2 vault (/var/lib/anchord/keys.luks), symlinked from
# <ANCHORD_HOME>/keyring-test → /opt/anchord/secure/keyring-test. Wait for
# the vault mount so the symlink resolves at boot. Never store the keyring
# in plaintext on an unencrypted disk (trust-layer.md W7, 2026-07-18).
RequiresMountsFor=/opt/anchord/secure

[Service]
Type=simple
User=b
Environment=ANCHORD_ADDR=:9110
Environment=ANCHORD_BIN=/home/b/go/bin/dreamtreed
Environment=ANCHORD_KEY=anchor
Environment=ANCHORD_KEYRING=test
Environment=ANCHORD_HOME=/home/b/.dreamtree
Environment=ANCHORD_CHAIN_ID=dreamtree
Environment=ANCHORD_NODE=tcp://localhost:16657
Environment=ANCHORD_FEES=0uphoton
# ANCHORD_TOKEN is supplied via an EnvironmentFile (from infra/.env), not committed:
EnvironmentFile=/etc/dreamtree/anchord.env
ExecStart=/home/b/go/bin/anchord
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
